# artnode-rpi/install/artnode-rpi.service
[Unit]
Description=artnode RPi — Art-Net LED controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/artnode-rpi
Restart=on-failure
RestartSec=5

# Root required for rpi-ws281x PWM/DMA hardware access.
# If you prefer not to run as root, grant CAP_SYS_RAWIO instead:
#   sudo setcap cap_sys_rawio+ep /usr/local/bin/artnode-rpi
# and uncomment the lines below:
# User=artnode
# Group=artnode
# AmbientCapabilities=CAP_SYS_RAWIO

# Config lives here; create the directory before first run:
#   sudo mkdir -p /etc/artnode
RuntimeDirectory=artnode
StateDirectory=artnode

# Limit core dumps and resource usage
LimitNOFILE=1024
LimitCORE=0

StandardOutput=journal
StandardError=journal
SyslogIdentifier=artnode-rpi

[Install]
WantedBy=multi-user.target
